# contrib/fsbackup.cron
#!/bin/sh
# Backup planner running from crontab.
#
#
# Example line for crontab: 18 4 * * * /opt/fsbackup/fs_backup.cron
#

letter="`/opt/fsbackup/create_backup.sh` "

if [ "$letter" ]; then 
    echo  "$letter " | mail -aFrom:"FROM NAME<from_example@example.com>" -s "Backup Report: $(hostname), $(hostname -I | awk '{print $1}')" to_example@example.com
fi
